# contrib/debian/prerm
#!/bin/sh
set -e
systemctl stop nomadterm 2>/dev/null || true
systemctl disable nomadterm 2>/dev/null || true
